# web/crontab
# Run cron_backup.php every day at 02:00 AM
0 2 * * * root /usr/local/bin/php /var/www/html/cron_backup.php >> /var/log/cron.log 2>&1
# An empty line is required at the end of this file for a valid cron file.
